# cron.hourly/srs-update
#!/usr/bin/env sh
echo -ne "UPDATE
BYE" | nc localhost 9876
